# recipes-support/emmc-installer/files/emmc-installer.service
[Unit]
Description=eMMC installer

[Service]
Type=oneshot
ExecStart=/usr/bin/emmc_installer.sh

[Install]
WantedBy=multi-user.target
